# server/etc/crontab
 */10 *  *  *  * ubuntu php /var/www/html/cli/comment-notify.php
 0   */5 *  *  * ubuntu php /var/www/html/cli/social-post.php
